# resources/docker/etc/supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisord.sock

[supervisord]
nodaemon=true
loglevel=debug
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.log

[supervisorctl]
serverurl=unix:///tmp/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:php-fpm]
command=/usr/sbin/php-fpm83 -F
autostart=true
autorestart=true
stderr_logfile=/var/log/php-fpm/stderr.log
stdout_logfile=/var/log/php-fpm/stdout.log

[program:nginx]
command=nginx
autorestart=true

# Leave it here until there will be a separate worker service
[program:queue-worker]
command=php /var/www/meal-to-feel-api/artisan queue:work database --sleep=3 --tries=2 --memory=500
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=3
stderr_logfile=/var/log/queue-worker/stderr.log
stdout_logfile=/var/log/queue-worker/stdout.log
